# x-c1-fan.service
[Unit]
Description=Run fan script xc-1

[Service]
Type=simple
ExecStart=/usr/local/bin/fan.py
Restart=always

[Install]
WantedBy=multi-user.target
